# amd64/rootfs/etc/services.d/udevd/run
#!/usr/bin/execlineb -P
with-contenv
background -d {
    s6-sleep 3
    udevadm trigger
}
udevd
